# deploy/temm1e.service
[Unit]
Description=TEMM1E AI Agent Runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=temm1e
Group=temm1e
WorkingDirectory=/var/lib/temm1e
EnvironmentFile=/etc/temm1e/env
ExecStart=/usr/local/bin/temm1e start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
